# timenowd/timenow.conf
start on started hildon-desktop
stop on stopped hal
#stop on stopped dbus
respawn
respawn limit 15 3

script
    exec su user -c "run-standalone.sh /opt/timenowd/bin/timenowd"

end script
